# src/etc/supervisord.d/vncserver.conf
[program:vncserver]
command=/bin/bash -c "/bin/vncserver-ctl.sh"
autostart=true
autorestart=true
user=user
redirect_stderr=true




